# platforms/osx/com.tonyseek.rsocks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tonyseek.rsocks</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/rsocks</string>
        <string>--config</string>
        <string>/usr/local/etc/rsocks/config.toml</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <!--
    <key>StandardOutPath</key>
    <string>/Users/YOUR_NAME/Library/Logs/rsocks.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_NAME/Library/Logs/rsocks.errors.log</string>
    -->
</dict>
</plist>
