# platforms/linux-systemd/rsocks.service
[Unit]
Description=A SOCKS reverse proxy server
After=network.target

[Service]
User=nobody
Group=nobody
WorkingDirectory=/opt/rsocks
ExecStart=/opt/rsocks/bin/rsocks --config=/etc/rsocks/config.toml

[Install]
WantedBy=multi-user.target
